# idata-ai-platform/scripts/systemd/idata-backend.service
[Unit]
Description=iData AI Platform Backend Services
After=network.target mysql.service redis.service rabbitmq-server.service
Wants=mysql.service redis.service rabbitmq-server.service

[Service]
Type=forking
User=xiaobofei
Group=xiaobofei
WorkingDirectory=/home/xiaobofei/apps/idata-ai-platform

# 启动脚本
ExecStart=/home/xiaobofei/apps/idata-ai-platform/scripts/deploy/start-services.sh all

# 停止脚本
ExecStop=/home/xiaobofei/apps/idata-ai-platform/scripts/deploy/stop-services.sh all

# 重启策略
Restart=on-failure
RestartSec=10

# 环境变量
Environment="JAVA_HOME=/usr/lib/jvm/java-17-openjdk-amd64"
Environment="PATH=/usr/lib/jvm/java-17-openjdk-amd64/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"

# 日志
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
